# config/systemd/monolith-backup.timer
[Unit]
Description=Monolith OS Daily Backup Timer

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
